# system_files/set-hdparm.service
[Unit]
Description=Set hdparm B255 on sda when booting

[Service]
ExecStart=/usr/sbin/hdparm -B255 /dev/sda

[Install]
WantedBy=multi-user.target 
